# deploy/ids-dashboard.service
[Unit]
Description=IDS Dashboard Service
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/ids-dashboard
ExecStart=/usr/bin/python3 -m ids.dashboard.main
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
